# lib/temp/javax-usb/rpm-files/etc/profile.d/javax-usb.sh

if test -n "$CLASSPATH" ; then
  CLASSPATH="/opt/javax-usb/lib/jsr80.jar:${CLASSPATH}"
else
  CLASSPATH="/opt/javax-usb/lib/jsr80.jar"
fi

export CLASSPATH
