# lib/temp/javax-usb/rpm-files/etc/profile.d/javax-usb.csh

if ( ${?CLASSPATH} ) then
  setenv CLASSPATH "/opt/javax-usb/lib/jsr80.jar:${CLASSPATH}"
else
  setenv CLASSPATH "/opt/javax-usb/lib/jsr80.jar"
fi
