# lib/temp/javax-usb-ri/rpm-files/etc/profile.d/javax-usb-ri.sh

if test -n "$CLASSPATH" ; then
  CLASSPATH="/opt/javax-usb/lib/jsr80_ri.jar:${CLASSPATH}"
else
  CLASSPATH="/opt/javax-usb/lib/jsr80_ri.jar"
fi

export CLASSPATH
